# systemd/karya.service
[Unit]
Description=karya — offline autonomous agent
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/karya
ExecStart=/home/pi/.local/bin/karya start --goals /home/pi/karya/config/goals.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=KARYA_HOME=/home/pi/.karya

[Install]
WantedBy=multi-user.target
